# deploy/dspm-bridge.service
[Unit]
Description=DSPM archive bridge (norns audience controller)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# --- edit these three for your Pi ---
User=pi
WorkingDirectory=/home/pi/dspm-archive/bridge
EnvironmentFile=/home/pi/dspm-archive/deploy/dspm-bridge.env
# node path: `which node` on the Pi (often /usr/bin/node, or an nvm path)
ExecStart=/usr/bin/node bridge-server.js
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
